# bunny-ai.service
[Unit]
Description=Jessica's Crabby Editor - Literary Analysis AI
After=network.target
Wants=network.target

[Service]
Type=simple
User=sean
Group=sean
WorkingDirectory=/Users/sean/src/bunnyAI
Environment=PATH=/Users/sean/src/bunnyAI/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/Users/sean/src/bunnyAI
ExecStart=/Users/sean/src/bunnyAI/venv/bin/python /Users/sean/src/bunnyAI/start_web.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/Users/sean/src/bunnyAI

[Install]
WantedBy=multi-user.target
